# files/etc/crontabs/root
0 4 * * * /usr/bin/wifi-auto-channel.sh
