# files/etc/uci-defaults/99-custom
#!/bin/sh

# =============================
# 5G 主 WiFi（高速 HE160）
# =============================
uci delete wireless.@wifi-iface[0]  # 删除默认接口
uci add wireless wifi-iface
uci set wireless.@wifi-iface[-1].device='radio1'
uci set wireless.@wifi-iface[-1].mode='ap'
uci set wireless.@wifi-iface[-1].ssid='AAAA_5G'
uci set wireless.@wifi-iface[-1].encryption='sae-mixed'
uci set wireless.@wifi-iface[-1].key='12345678'
uci set wireless.radio1.htmode='HE160'
uci set wireless.radio1.country='CN'
uci set wireless.radio1.vht_highest=1
uci set wireless.radio1.txpower=20
uci set wireless.radio1.ampdu=1
uci set wireless.radio1.wmm=1

# =============================
# 2.4G IoT WiFi（门锁专用）
# =============================
uci add wireless wifi-iface
uci set wireless.@wifi-iface[-1].device='radio0'
uci set wireless.@wifi-iface[-1].mode='ap'
uci set wireless.@wifi-iface[-1].ssid='AAAA'
uci set wireless.@wifi-iface[-1].encryption='psk2'
uci set wireless.@wifi-iface[-1].key='12345678'
uci set wireless.@wifi-iface[-1].ieee80211w='0'
uci set wireless.@wifi-iface[-1].disassoc_low_ack='0'
uci set wireless.radio0.country='CN'
uci set wireless.radio0.htmode='HT20'
uci set wireless.radio0.hwmode='11g'
uci set wireless.radio0.txpower=18
uci set wireless.radio0.ampdu=0
uci set wireless.radio0.wmm=1

# 提交配置
uci commit wireless

exit 0
